# deploy/docker/com.claw-empire.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claw-empire.tunnel</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/cloudflared</string>
        <string>tunnel</string>
        <string>run</string>
        <string>--token-file</string>
        <string>/Users/KangHyeonSeok/Documents/GitHub/claw-empire/deploy/docker/runtime/cloudflared-token.txt</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/KangHyeonSeok/Documents/GitHub/claw-empire/logs/cloudflared-tunnel.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/KangHyeonSeok/Documents/GitHub/claw-empire/logs/cloudflared-tunnel.error.log</string>
</dict>
</plist>